# deploy/quant-sync-a.timer
[Unit]
Description=AI Quant Data Subsystem — A股同步定时器（沪深，16:30 北京时间）
Requires=quant-sync.service

[Timer]
# A股收盘后 16:30 北京时间触发
# TimeZone=Asia/Shanghai 显式声明时区，避免 systemd 服务器为 UTC 时差 8 小时
OnCalendar=Mon-Fri 16:30:00
TimeZone=Asia/Shanghai
# 允许 ±5 分钟的随机延迟，错峰启动
RandomizedDelaySec=300
# 错过触发时（如机器关机）重启后立即补跑一次
Persistent=true
Unit=quant-sync.service

[Install]
WantedBy=timers.target
