# deploy/quant-sync-us.timer
[Unit]
Description=AI Quant Data Subsystem — 美股同步定时器（US，07:00 北京时间次日）
Requires=quant-sync.service

[Timer]
# 美股收盘（东部时间 16:00 / 北京时间次日 04:00(夏令) 或 05:00(冬令)）后，07:00 北京时间触发
# 07:00 为固定安全时间，无需区分冬夏令时
# TimeZone=Asia/Shanghai 显式声明时区，避免 systemd 服务器为 UTC 时差 8 小时
OnCalendar=Tue-Sat 07:00:00
TimeZone=Asia/Shanghai
# 允许 ±5 分钟的随机延迟，错峰启动
RandomizedDelaySec=300
# 错过触发时（如机器关机）重启后立即补跑一次
Persistent=true
Unit=quant-sync.service

[Install]
WantedBy=timers.target
